# laravel-backend/docker/conf/supervisord.websocket.conf
[supervisord]
nodaemon=true
[program:app-websoket]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan websockets:serve --port=6005
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


